# master_orchestrator/launchd/com.trading.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Job Label -->
    <key>Label</key>
    <string>com.trading.daily</string>

    <!-- Program to run (with drive wait wrapper) -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/henry/miniconda3/envs/trading/bin/python</string>
        <string>/Users/henry/Desktop/Programming &amp; Networking/Personal Projects/Integrated Trading Agent/master_orchestrator/wait_for_drive.py</string>
        <string>--daily</string>
        <string>--max-wait-hours</string>
        <string>2</string>
    </array>

    <!-- Working Directory -->
    <key>WorkingDirectory</key>
    <string>/Users/henry/Desktop/Programming &amp; Networking/Personal Projects/Integrated Trading Agent</string>

    <!-- Run every day at 5:00 PM ET (after market close at 4:00 PM) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>17</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Standard Output -->
    <key>StandardOutPath</key>
    <string>/tmp/trading_daily.log</string>

    <!-- Standard Error -->
    <key>StandardErrorPath</key>
    <string>/tmp/trading_daily_error.log</string>

    <!-- Run at load (for testing) -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Environment Variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/henry/miniconda3/envs/trading/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>TRADING_EMAIL_SENDER</key>
        <string>henry.vianna123@gmail.com</string>
        <key>TRADING_EMAIL_PASSWORD</key>
        <string>bhfy yrsc sbyy efpj</string>
    </dict>
</dict>
</plist>
